# conf/named.conf
[program:bind9]
command=/usr/sbin/named -f -d 3 
stdout_events_enabled=true
stderr_events_enabled=true
